# deploy/systemd/qemubox-nexuserofs-snapshotter.service
[Unit]
Description=nexuserofs external snapshotter for containerd
Documentation=https://github.com/aledbf/nexuserofs
After=network.target
Before=qemubox-containerd.service

[Service]
ExecSearchPath=/usr/share/qemubox/bin
ExecSearchPath=/usr/local/sbin
ExecSearchPath=/usr/local/bin
ExecSearchPath=/usr/sbin
ExecSearchPath=/usr/bin
ExecSearchPath=/sbin
ExecSearchPath=/bin
ExecStartPre=-/bin/mkdir -p /run/qemubox/nexuserofs-snapshotter
ExecStartPre=-/bin/mkdir -p /var/lib/qemubox/nexuserofs-snapshotter
ExecStart=/usr/share/qemubox/bin/nexuserofs-snapshotter \
  --root /var/lib/qemubox/nexuserofs-snapshotter \
  --address /run/qemubox/nexuserofs-snapshotter/snapshotter.sock \
  --containerd-address /var/run/qemubox/containerd.sock \
  --containerd-namespace default \
  --log-level info \
  --set-immutable \
  --mkfs-options="-zlz4hc,12" \
  --mkfs-options="-C65536"

Type=simple
Restart=always
RestartSec=5

# Wait for the socket to be ready before declaring service as started
ExecStartPost=/bin/sh -c 'for i in $(seq 1 30); do [ -S /run/qemubox/nexuserofs-snapshotter/snapshotter.sock ] && exit 0; sleep 0.1; done; exit 1'

[Install]
WantedBy=multi-user.target
RequiredBy=qemubox-containerd.service
