# deploy/systemd/qemubox-nexus-erofs-snapshotter.service
[Unit]
Description=nexus-erofs external snapshotter for containerd
Documentation=https://github.com/aledbf/nexus-erofs
After=network.target
Before=qemubox-containerd.service

[Service]
ExecSearchPath=/usr/share/qemubox/bin
ExecSearchPath=/usr/local/sbin
ExecSearchPath=/usr/local/bin
ExecSearchPath=/usr/sbin
ExecSearchPath=/usr/bin
ExecSearchPath=/sbin
ExecSearchPath=/bin
ExecStartPre=-/bin/mkdir -p /run/qemubox/nexus-erofs-snapshotter
ExecStartPre=-/bin/mkdir -p /var/lib/qemubox/nexus-erofs-snapshotter
ExecStart=/usr/share/qemubox/bin/nexus-erofs-snapshotter \
  --root /var/lib/qemubox/nexus-erofs-snapshotter \
  --address /run/qemubox/nexus-erofs-snapshotter/snapshotter.sock \
  --containerd-address /var/run/qemubox/containerd.sock \
  --containerd-namespace default \
  --log-level debug \
  --set-immutable

Type=simple
Restart=always
RestartSec=5

# Wait for the socket to be ready before declaring service as started
ExecStartPost=/bin/sh -c 'for i in $(seq 1 30); do [ -S /run/qemubox/nexus-erofs-snapshotter/snapshotter.sock ] && exit 0; sleep 0.1; done; exit 1'

[Install]
WantedBy=multi-user.target
RequiredBy=qemubox-containerd.service
