# images/sandbox/rootfs/etc/udev/rules.d/disable-udev-settle.rules
# Minimal udev rules for VM - skip unnecessary processing
ACTION=="add", SUBSYSTEM=="net", RUN+="/bin/true"
ACTION=="remove", SUBSYSTEM=="net", RUN+="/bin/true"
